# rforms.service
[Unit]
Description=RForms Django Application
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/rforms
Environment="PATH=/home/ec2-user/rforms/venv/bin"
ExecStart=/home/ec2-user/rforms/venv/bin/python manage.py runserver 0.0.0.0:8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
